# monitor_process@.service
[Unit]
Description=Monitor for process %I

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor_process.sh %I
